# skills/memory-reflect/memory-daemon.service
[Unit]
Description=Total Recall Daemon
After=network.target
# Запускаться после того как Neo4j и Qdrant доступны
# (если они тоже сервисы — добавить After=neo4j.service)

[Service]
Type=simple
User=ironman
WorkingDirectory=/home/ironman/.openclaw/skills/memory-reflect
ExecStart=/home/ironman/.openclaw/skills/memory-reflect/.venv/bin/python memory-daemon.py
Restart=on-failure
RestartSec=30
# Демон не падает сам — но если упал, подождать 30с перед рестартом

# Логи через journald (дополнительно к файловым логам из store.py)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=total-recall-daemon

# .env должен лежать в WorkingDirectory
EnvironmentFile=/home/ironman/.openclaw/skills/memory-reflect/.env

[Install]
WantedBy=multi-user.target
